# home/pi/.config/autostart/gpstrace.desktop
[Desktop Entry]
Type=Application
Name=GPS Trace
Exec=/home/pi/gps/startgpstrace.sh
StartupNotify=false
