# home/pi/.config/autostart/navit.desktop
[Desktop Entry]
Type=Application
Name=Navit
Exec=/usr/bin/navit
StartupNotify=false
